# CivitaiFirehose/civitai-firehose.service
[Unit]
Description=Civitai Firehose
After=network.target

[Service]
Type=exec
User=jackkendall

WorkingDirectory=/opt/CivitaiFirehose
ExecStart=/usr/bin/dotnet /opt/CivitaiFirehose/CivitaiFirehose.dll

Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5124

Restart=always
RestartSec=10
KillSignal=SIGINT

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

CPUQuota=30%
MemoryMax=512M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=civitai-firehose

[Install]
WantedBy=multi-user.target